# data/com.system76.Scheduler.service
[Unit]
Description=Automatically configure CPU scheduler for responsiveness on AC

[Service]
ExecStartPre=/bin/sleep 10
ExecStart=/usr/bin/system76-scheduler daemon
ExecReload=/usr/bin/system76-scheduler daemon reload
Type=dbus
BusName=com.system76.Scheduler

[Install]
WantedBy=multi-user.target
